# couchpotato.ini
[program:CouchPotatoServer]
command=/CouchPotato/CouchPotato.py --data_dir=/config
