# couchpotato.ini
[program:CouchPotatoServer]
user=
command=/CouchPotatoServer/CouchPotato.py --data_dir=/config
