# couchpotato.ini
[program:CouchPotatoServer]
command=/CouchPotatoServer/CouchPotato.py --data_dir=/config
